# crontab
0 */1 * * * curl -X POST http://localhost:3000/api/admin/populate

